# config/runit/postfix/run
#!/bin/sh
exec /usr/sbin/postfix start-fg
